# bib-base/files/rootdir/etc/profile.d/editor.sh
export EDITOR=/usr/bin/lvim
export VISUAL=/usr/bin/lvim
