# bib-base/files/rootdir/etc/profile.d/editor.sh
export EDITOR=$HOME/.local/bin/lvim
export VISUAL=$HOME/.local/bin/lvim
