# bib-base/files/rootdir/etc/profile.d/login.sh
echo "Starting SE bastille installer.."
doas ./SE_bastille-installer/SE_bastille-installer
